# private_dot_config/autostart/autostart.desktop
[Desktop Entry]
Type=Application
Name=Startup script
Terminal=false
Exec=sh -c "${HOME}/.local/bin/autostart"
